# DEBIAN/prerm
#!/bin/bash

systemctl stop debian-demo
systemctl disable debian-demo